# deploy/systemd/deepseek-native-agent.service
[Unit]
Description=DeepSeek Native Agent Docker Compose gateway
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/deepseek-native-agent/deploy
ExecStart=/usr/bin/docker compose --env-file %h/deepseek-native-agent/deploy/.env up -d
ExecStop=/usr/bin/docker compose --env-file %h/deepseek-native-agent/deploy/.env down
TimeoutStartSec=180
TimeoutStopSec=120

[Install]
WantedBy=default.target
